# crontabs/root
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.88Yl0d/crontab installed on Wed Mar 31 14:58:47 2021)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
# Edit this file to introduce tasks to be run by cron.
#
# Each task to run has to be defined through a single line
# indicating with different fields when the task will be run
# and what command to run for the task
#
# To define the time you can provide concrete values for
# minute (m), hour (h), day of month (dom), month (mon),
# and day of week (dow) or use '*' in these fields (for 'any').
#
# Notice that tasks will be started based on the cron's system
# daemon's notion of time and timezones.
#
# Output of the crontab jobs (including errors) is sent through
# email to the user the crontab file belongs to (unless redirected).
#
# For example, you can run a backup of all your user accounts
# at 5 a.m every week with:
# 0 5 * * 1 tar -zcf /var/backups/home.tgz /home/
#
# For more information see the manual pages of crontab(5) and cron(8)
#
# m h  dom mon dow   command
#m h dom mon dow command
#每分钟检查wifi连接 无法ping通192.168.50.1则重连
#*/1 * * * * /var/opt/checkWifiConnect.sh

#每天6点重新绑定路由器arp 用于远程唤醒
#* 6 * * * /var/opt/arpInRouter.sh

#每半小时备份我的世界
#*/30 * * * * /var/opt/bedrock-server-backup/bedrock-server-backup.sh

# 每30分钟拉一次流
*/30 * * * * /opt/surveillance-video.sh

#创建前一日日期文件夹 并将监控视频移至该文件夹
01 00 * * * mkdir /d/ftp/haikang/$(date -d "-1 days" +\%Y\%m\%d) && mv /d/ftp/haikang/*$(date -d "-1 days" +\%Y\%m\%d)*.* /d/ftp/haikang/$(date -d "-1 days" +\%Y\%m\%d)

# 删除40天前的录像记录
00 00 * * * rm -r /d/ftp/haikang/$(date +\%Y\%m\%d -d "-40 days")

# 测试crontab
* * * * * echo $(date) >> /tmp/cron_root_lob.txt
